# roles/docker-elasticsearch/templates/systemd-container.service.j2
[Unit]
Description={{ container_description }}
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker stop {{ container_name }}
ExecStartPre=-/usr/bin/docker rm {{ container_name }}
{% if dockerhub == true %}
ExecStartPre=-/usr/bin/docker pull {{ full_image }}
{% endif %}
ExecStart=/usr/bin/docker run --rm -t \
{% for port in container_ports %}
    -p {{ port }} \
{% endfor %}
{% for env in env_variables %}
    -e {{ env }} \
{% endfor %}
{% for volume in container_volumes %}
    -v {{ volume }} \
{% endfor %}
{% for option in container_extra_options %}
    {{ option }} \
{% endfor %}
    --name {{ container_name }} {{ full_image }} {{ container_cmd }}
ExecStop=-/usr/bin/docker stop -t 3 {{ container_name }}
ExecStop=-/usr/bin/docker rm {{ container_name }}
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target